# .config/aptik/initial-files-cron.list
/etc/cron.daily
/etc/cron.daily/logrotate
/etc/cron.daily
/etc/cron.daily/mlocate
/etc/cron.daily
/etc/cron.daily/dpkg
/etc/cron.d
/etc/cron.daily
/etc/cron.daily/popularity-contest
/etc/cron.d
/etc/cron.d/anacron
/etc/cron.daily
/etc/cron.daily/0anacron
/etc/cron.monthly
/etc/cron.monthly/0anacron
/etc/cron.weekly
/etc/cron.weekly/0anacron
/etc/cron.daily
/etc/cron.daily/update-notifier-common
/etc/cron.weekly
/etc/cron.weekly/update-notifier-common
/etc/cron.daily
/etc/cron.daily/apt-compat
/etc/cron.d
/etc/cron.d/.placeholder
/etc/cron.daily
/etc/cron.daily/.placeholder
/etc/cron.hourly
/etc/cron.hourly/.placeholder
/etc/cron.monthly
/etc/cron.monthly/.placeholder
/etc/cron.weekly
/etc/cron.weekly/.placeholder
/etc/crontab
/etc/cron.daily
/etc/cron.daily/apport
/etc/cron.daily
/etc/cron.daily/passwd
/etc/cron.daily
/etc/cron.daily/bsdmainutils
/etc/cron.daily
/etc/cron.daily/man-db
/etc/cron.weekly
/etc/cron.weekly/man-db
/etc/cron.daily
/etc/cron.daily/cracklib-runtime
/etc/cron.daily
/etc/cron.daily/ubuntu-advantage-tools
